# cookbooks/stl-rancher/templates/etc/systemd/system/docker.service.d/http-proxy.conf.erb
[Service]
Environment="HTTP_PROXY=http://proxy.tabajara.intranet:3130/"
Environment="HTTPS_PROXY=http://proxy.tabajara.intranet:3130/"
Environment="NO_PROXY=.tabajara.local, .tabajara.intranet, 10.0.0.0/8"
