# cookbooks/stl-infra-scripts/files/cron.d/backup-git-infra
# Backup do git de infra
20 2 * * * root /bin/sleep `/usr/bin/expr $RANDOM \% 60`;/usr/local/bin/backup-git-infra.sh >/dev/null 2>&1

